# root/etc/s6-overlay/s6-rc.d/init-hashserv-data/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# permissions
lsiown abc:abc \
	/data \
	/data/*
